# debian/postinst
#!/bin/sh
set -e
case "$1" in
    configure)
        # Create system user if it does not exist
        if ! id -u waveshare-ups > /dev/null 2>&1; then
            adduser --system --no-create-home --group \
                    --home /nonexistent \
                    --shell /usr/sbin/nologin \
                    waveshare-ups
        fi

        # Grant I2C access
        if getent group i2c > /dev/null 2>&1; then
            usermod -aG i2c waveshare-ups
        fi

        # Lock down config
        chown root:waveshare-ups /etc/waveshare-ups/config.ini
        chmod 640 /etc/waveshare-ups/config.ini
        ;;
esac
#DEBHELPER#

echo ""
echo "waveshare-ups installed."
echo "  Edit config : /etc/waveshare-ups/config.ini"
echo "                Verify device.addr matches your board (default 0x43)"
echo "                Run 'i2cdetect -y 1' to find the address"
echo "  Enable      : systemctl enable --now waveshare-ups"
echo "  Stream data : nc 127.0.0.1 8564"
echo "  Status      : echo status | nc 127.0.0.1 8565"
echo ""
exit 0
