# debian/prerm
#!/bin/sh
set -e
case "$1" in
    remove|purge)
        if systemctl is-active --quiet waveshare-ups.service 2>/dev/null; then
            systemctl stop waveshare-ups.service || true
        fi
        if systemctl is-enabled --quiet waveshare-ups.service 2>/dev/null; then
            systemctl disable waveshare-ups.service || true
        fi
        ;;
esac
#DEBHELPER#
exit 0
